# deploy/nix/systemd/windmill-extra.service
# SPDX-License-Identifier: copyleft-next-0.3.1
# Windmill extra (the LSP gateway for editor intellisense), run from nix under
# systemd --user. Mirrors deploy/podman/systemd/windmill-extra.container. The
# Tornado gateway listens on PORT and spawns a language server per session;
# caddy routes the /ws* paths here. WINDMILL_BASE_URL lets it pull the instance
# ruff config from the server.
#
# Binary out-link:
#   nix build .#windmill-extra --out-link %S/windmill/sw/windmill-extra
# Override PORT with `systemctl --user edit` or %E/windmill/windmill-extra.env.
[Unit]
Description=Windmill LSP Gateway
Requires=windmill.service
After=windmill.service

[Service]
Environment=PORT=3001
Environment=WINDMILL_BASE_URL=http://127.0.0.1:8002
EnvironmentFile=-%E/windmill/windmill-extra.env
ExecStart=%S/windmill/sw/windmill-extra/bin/windmill-extra
Restart=always

[Install]
WantedBy=default.target
