# deploy/nix/systemd/windmill-worker-vmrun@.service
# SPDX-License-Identifier: copyleft-next-0.3.1
# Windmill VM long-poll worker (group vm, tag vm-run), a systemd template:
# windmill-worker-vmrun@1, @2, ... Same host integration as the vm worker, but
# pulls only the vm-run tag: the long-lived fstests wait poll, split from the vm
# control tag so a multi-hour run never starves a quick lifecycle op. The
# instance count is the concurrent-test-run cap. Mirrors the vm-run side of
# deploy/podman/windmill-worker-vm.container.tmpl.
#
# Not enabled by default: it needs the workbench provisioned and vhost_vsock
# loaded. Override the build-area with `systemctl --user edit` or
# %E/windmill-nix/windmill-worker-vmrun.env, then
# `systemctl --user enable --now windmill-worker-vmrun@1`.
[Unit]
Description=Windmill vm-run worker %i (group vm; long fstests poll) (nix)
Requires=windmill-db.service
After=windmill-db.service

[Service]
Environment=MODE=worker
Environment=NUM_WORKERS=1
Environment=WORKER_GROUP=vm
Environment=WORKER_TAGS=vm-run
Environment=WORKER_INDEX=vmrun%i
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
Environment=WORKBENCH_DIR=%S/windmill-nix/workbench
Environment=SYSTEM_DIR=%S/windmill-nix/workbench/system
Environment=WORKERS_DIR=%S/windmill-nix/workbench/workers
Environment=VENDOR_DIR=%S/windmill-nix/vendor
Environment=WHITELIST_ENVS=DBUS_SESSION_BUS_ADDRESS,XDG_RUNTIME_DIR,HOME,XDG_CONFIG_HOME,XDG_STATE_HOME,WORKER_INDEX,WORKBENCH_DIR,SYSTEM_DIR,WORKERS_DIR,VENDOR_DIR
EnvironmentFile=%S/windmill-nix/env/database.env
EnvironmentFile=-%E/windmill-nix/windmill-worker-vmrun.env
ExecStart=%S/windmill-nix/sw/windmill/bin/windmill
Restart=always

[Install]
WantedBy=default.target
