# deploy/nix/systemd/windmill-worker-vm@.service
# SPDX-License-Identifier: copyleft-next-0.3.1
# Windmill VM-orchestration worker (group vm, tag vm), a systemd template:
# windmill-worker-vm@1, @2, ... run the f/qsu/* control and lifecycle steps that
# render qsu's systemd units and drive the host systemd --user to boot QEMU VMs.
# Mirrors deploy/podman/windmill-worker-vm.container.tmpl for the vm tag.
#
# Under systemd --user the podman template's whole host-integration apparatus is
# free: no container means the user bus, /dev/vsock, and the host
# ~/.config/systemd + ~/.local/state/qemu-system search paths are already this
# user's, so the bind-mounts, AddDevice and the AF_VSOCK seccomp profile all
# drop. The job still needs the bus address and the XDG dirs passed through to
# the steps (WHITELIST_ENVS) and the build area below.
#
# Not enabled by default: it needs the workbench provisioned (the System bare,
# the peer ssh key) and the vhost_vsock module loaded. Point the build-area at a
# real workbench (`systemctl --user edit` or %E/windmill-nix/windmill-worker-vm.env)
# and `systemctl --user enable --now windmill-worker-vm@1`.
[Unit]
Description=Windmill vm worker %i (group vm; QEMU/systemd VM orchestration) (nix)
Requires=windmill-db.service
After=windmill-db.service

[Service]
Environment=MODE=worker
Environment=NUM_WORKERS=1
Environment=WORKER_GROUP=vm
Environment=WORKER_TAGS=vm
Environment=WORKER_INDEX=vm%i
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
Environment=WORKBENCH_DIR=%S/windmill-nix/workbench
Environment=SYSTEM_DIR=%S/windmill-nix/workbench/system
Environment=WORKERS_DIR=%S/windmill-nix/workbench/workers
Environment=VENDOR_DIR=%S/windmill-nix/vendor
Environment=WHITELIST_ENVS=DBUS_SESSION_BUS_ADDRESS,XDG_RUNTIME_DIR,HOME,XDG_CONFIG_HOME,XDG_STATE_HOME,WORKER_INDEX,WORKBENCH_DIR,SYSTEM_DIR,WORKERS_DIR,VENDOR_DIR
EnvironmentFile=%S/windmill-nix/env/database.env
EnvironmentFile=-%E/windmill-nix/windmill-worker-vm.env
ExecStart=%S/windmill-nix/sw/windmill/bin/windmill
Restart=always

[Install]
WantedBy=default.target
