# deploy/podman/windmill-worker.container.tmpl
# SPDX-License-Identifier: copyleft-next-0.3.1
# Template: install.sh renders one copy per worker, replacing @INDEX@.
[Unit]
Description=Windmill worker @INDEX@ (group default; builds)
Requires=windmill-db.service
After=windmill-db.service
[Container]
ContainerName=windmill-worker-@INDEX@
Image=ghcr.io/windmill-labs/windmill:main
Network=windmill.network
Environment=MODE=worker
Environment=NUM_WORKERS=1
Environment=WORKER_GROUP=default
Environment=DATABASE_URL=postgres://postgres:changeme@windmill-db:5432/windmill?sslmode=disable
# Drive the host/controller systemd --user (to spin up machines) over its socket.
Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
Environment=HOME=%h
Environment=WHITELIST_ENVS=DBUS_SESSION_BUS_ADDRESS,XDG_RUNTIME_DIR,WORKER_INDEX,WORKBENCH_DIR,SYSTEM_DIR,WORKERS_DIR,VENDOR_DIR
Environment=WORKER_INDEX=@SELFDIR@
Environment=WORKBENCH_DIR=@WORKBENCH_DIR@
Environment=SYSTEM_DIR=@SYSTEM_DIR@
Environment=WORKERS_DIR=@WORKERS_DIR@
Environment=VENDOR_DIR=@VENDOR_DIR@
Volume=%t/bus:%t/bus
# /nix store, bind-mounted for building. The git mirror the worktrees borrow
# objects from lives in the System workbench (@SYSTEM_DIR@/mirror, mounted
# below), so no separate /mirror mount is needed.
Volume=/nix:/nix
Volume=%C/windmill/cache:/tmp/windmill/cache
# Shared logs volume (see windmill.container): workers write large job logs here
# so the server can serve them for download.
Volume=%C/windmill/logs:/tmp/windmill/logs
# Build area at IDENTICAL host paths so git worktree links resolve. By default the
# whole Workbench is mounted (the System workbench, every worker sandbox, the shared
# caches, and the developer worktree-groups), so a worker resolves any of them and
# the worktree-group initializer can lay developer groups down host-visibly. The
# pinned vendored deps are a sibling of the Workbench, mounted read-only (ADR-0006).
# When SYSTEM_DIR or WORKERS_DIR is relocated outside the Workbench, install.sh fills
# in a separate mount for it (otherwise the line is empty).
Volume=@WORKBENCH_DIR@:@WORKBENCH_DIR@:rw
@SYSTEM_MOUNT@
@WORKERS_MOUNT@
Volume=@VENDOR_DIR@:@VENDOR_DIR@:ro
[Service]
Restart=always
[Install]
WantedBy=default.target
